# .crontab
# ruby gem nightly updates
0 1 * * *  rvm get head;rvm reload
0 2 * * *  gem update -V
# brew nightly updates
0 3 * * *  brew update
0 4 * * *  brew upgrade
# npm
0 5 * * *  npm update npm -g
